# deployments/systemd-unit-files/clean-docker-images.service
[Unit]
Description="Clean up docker images"

[Service]
Type=simple
ExecStart=/usr/bin/docker system prune --all --force
User={CHANGE_ME}
